# arch/conf/crontab/prosdkr
* * * * * /home/prosdkr/.config/i3status/uptime.sh
